# config/systemd/user/minute-timer.target
##
## /etc/systemd/system/minute-timer.target
##
[Unit]
Description=Minute Timer Target
StopWhenUnneeded=yes
