# config/systemd/user/user-backup.target
##
## ~/.config/systemd/user/user-backup.target
##
[Unit]
Description=Back up user profile
StopWhenUnneeded=yes
